# libvirtd-desktop/usr/lib/systemd/system/virtstoraged.service.d/after-ldreload.conf
# Make sure we order after the linker reload
[Unit]
After=ldconfig.service
