# libvirtd-desktop/files/etc/systemd/system/virtlogd.service.d/override.conf
# SPDX-FileCopyrightText: 2025 Cody Istre
# SPDX-License-Identifier: MIT

[Service]
LogsDirectory=libvirt/qemu
ReadWritePaths=/var/log/libvirt
