# libvirtd-desktop/files/etc/systemd/system/virtlogd.service.d/override.conf
[Service]
LogsDirectory=libvirt/qemu
ReadWritePaths=/var/log/libvirt
