# libvirtd/usr/lib/systemd/system/virtlogd.service.d/override.conf
[Service]
ReadWritePaths=/var/log/libvirt/qemu
