# libvirtd-desktop/usr/lib/systemd/system/libvirt-dbus.service.d/reload-dbus.conf
[Service]
# We use '+' to run this specific command as root, 
# even if the rest of the service runs as a restricted user.
# We do this because the dbus policy file allowing this service to connect
# to the bus is injected after dbus already started.
ExecStartPre=+/usr/bin/busctl call org.freedesktop.DBus /org/freedesktop/DBus org.freedesktop.DBus ReloadConfig
