# libvirtd/usr/lib/systemd/system/virtstoraged.service.d/after-ldreload.conf
# Make sure we order after the linker service
[Unit]
After=ldconfig.service
